# src/templates/dataplane@.service
[Unit]
Description=Router data plane as pod

[Service]
ExecStart=/usr/local/sbin/setup-dataplane %i

[Install]
WantedBy=multi-user.target
